# etc/local.d/restore_brightness.start
#!/bin/sh
/home/akash/bin/set_brightness r
